# zapret-ebpf/files/etc/init.d/zapret-ebpf
#!/bin/sh /etc/rc.common
# SPDX-License-Identifier: MIT
# Запуск eBPF-фильтров Zapret на br-lan (WAN)

START=99
STOP=10
USE_PROCD=1

start() {
    echo "Starting zapret-ebpf on interface br-lan..."

    EBPF_DIR=/usr/lib/zapret/ebpf

    # Загружаем каждый BPF объект на XDP/TC
    for f in $EBPF_DIR/*.bpf.o; do
        fname=$(basename $f)
        echo "Attaching $fname to br-lan..."
        # Используем tc clsact qdisc для ingress
        tc qdisc add dev br-lan clsact 2>/dev/null || true
        tc filter add dev br-lan ingress bpf da obj $f sec prog 2>/dev/null || true
    done

    echo "zapret-ebpf started."
}

stop() {
    echo "Stopping zapret-ebpf on br-lan..."
    tc qdisc del dev br-lan clsact 2>/dev/null || true
    echo "zapret-ebpf stopped."
}

restart() {
    stop
    start
}
