# zapret-ebpf/files/etc/init.d/zapret-ebpf
#!/bin/sh /etc/rc.common
# Copyright (C) 2025
# zapret-ebpf init script
# Loads eBPF filters at boot

START=99
STOP=10
USE_PROCD=1

start() {
    echo "Loading zapret eBPF programs..."

    EBPF_DIR=/usr/lib/zapret/ebpf

    if [ ! -d "$EBPF_DIR" ]; then
        echo "Directory $EBPF_DIR not found!"
        return 1
    fi

    for prog in $EBPF_DIR/*.bpf.o; do
        name=$(basename $prog .bpf.o)
        echo "Loading eBPF program $name..."
        # Пример загрузки через tc (XDP):
        tc qdisc add dev eth0 clsact 2>/dev/null || true
        tc filter add dev eth0 ingress bpf obj $prog sec classifier || true
    done

    echo "zapret-ebpf loaded."
}

stop() {
    echo "Unloading zapret eBPF programs..."
    # tc filter delete можно расширить по мере необходимости
    tc filter delete dev br-lan ingress 2>/dev/null || true
    echo "zapret-ebpf unloaded."
}
